# systemd/mediarr-update.timer
[Unit]
Description=Run Mediarr update weekly

[Timer]
OnCalendar=Sun 02:00:00
AccuracySec=1h
Persistent=true

[Install]
WantedBy=timers.target
